# sleep_mode/sleep_mode.service
[Unit]
Description=low power mode service

[Service]
Type=oneshot
ExecStart=/usr/sbin/sleep_mode/sleep_mode.sh

[Install]
WantedBy=multi-user.target